# src/systemd-units/main.target
[Unit]
Description=Update vendor packages
Wants=pacman-update.target flatpak-update.service
StopWhenUnneeded=true
