# src/type/pacman/notifier/user/update/main.service
[Unit]
Description=Pacman upgrade progress notifier

[Service]
Type=exec
KillSignal=INT
ExecStart=/usr/lib/pacman-update_notifiers/user/update
